# Scripts/elysia.service
[Unit]
Description=Elysia Cognitive Daemon
After=network.target ollama.service

[Service]
Type=simple
User=%u
WorkingDirectory=/home/%u/Elysia
ExecStart=/usr/bin/python3 /home/%u/Elysia/elysia.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
